# admin/lib/systemd/system/dkim-rekey@.service
[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/dkim-rekey %i
