# admin/lib/systemd/system/dnssec-roll-ksk@.service
[Service]
Type=simple
User=named
Group=named
ExecStart=/usr/local/bin/dnssec-roll-ksk %i
