# admin/lib/systemd/system/dnssec-roll-ksk@.timer
[Timer]
OnCalendar=*-12-01 08:23:00
Persistent=true
AccuracySec=1s
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
