# pacman/lib/systemd/system/pacman-refresh.service
[Service]
Type = oneshot
ExecStart = pacman -Sy
