# minecraft/lib/systemd/system/backup-minecraft@.service
[Service]
User = %i
ExecStart = /usr/local/bin/backup-minecraft %i
