# minecraft/lib/systemd/system/backup-minecraft@.timer
[Timer]
AccuracySec = 1s
OnCalendar = 00..23:00..59/5:00

[Install]
WantedBy = timers.target
